# scripts/launchd/com.jchytech.pka-plane-autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jchytech.pka-plane-autostart</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jchavauxm5/PKA_JCH/bin/plane-autostart.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jchavauxm5/PKA_JCH</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>/Users/jchavauxm5/PKA_JCH/tmp/plane-autostart.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jchavauxm5/PKA_JCH/tmp/plane-autostart.launchd.log</string>
</dict>
</plist>
