# scripts/launchd/com.jchytech.pka-vault-maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jchytech.pka-vault-maintenance</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
    <string>/Users/jchavauxm5/PKA_JCH/scripts/pka_vault_maintenance.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jchavauxm5/PKA_JCH</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/jchavauxm5/PKA_JCH/scripts/logs/vault_maintenance.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jchavauxm5/PKA_JCH/scripts/logs/vault_maintenance.log</string>
</dict>
</plist>
